# chin.ini
[db]
sql_conn = mysql+mysqldb://chin:chin+1s@nazgrim.com/chin?charset=utf8

[webserver]
host = 0.0.0.0
port = 8443
secret_key = anything_yout

[scheduler]
heartbeat_sec = 5

[worker]
heartbeat_sec = 5

[uwsgi]
http = 0.0.0.0:8443
chdir = /home/suit/chin
wsgi-file = chin.py
callable = app 
processes = 4
threads = 2

[program:chin]
command=env PATH="/home/suit/dev/anaconda/bin" uwsgi /home/suit/chin/chin.ini
directory=/home/suit/chin
user=root
autostart=true
autorestart=true
stdout_logfile=/home/suit/chin/logs/uwsgi_supervisor.log      
die-on-term=true

